# home/pi/.config/autostart/vita.desktop
[Desktop Entry]
Type=Application
Name=Vita
Exec=/home/pi/autorun.sh