# etc/systemd/user/switch.service
[Service]
ExecStart=/home/pi/runswitch.sh
User=pi
StandardOutput=tty1