# etc/systemd/user/vita.service
[Service]
ExecStart=/home/pi/run.sh
User=pi
StandardOutput=tty1
