# etc/udev/rules.d/94-dswitch.rules
ACTION=="remove", SUBSYSTEM=="usb", RUN+="/home/pi/dswitch.sh"
