# etc/xdg/autostart/1stboot.desktop
[Desktop Entry]
Type=Application
Name=1stboot.desktop
Exec=./1stboot.sh
